# profile.d/dev-on-osx

LANG=en_US.UTF-8

export GIT_PAGER='less -FXR'

[ -f /usr/local/etc/bash_completion ] && . /usr/local/etc/bash_completion

# check the window size after each command and, if necessary,
# update the values of LINES and COLUMNS.
shopt -s checkwinsize

export CLICOLOR=1
export LSCOLORS=DxGxcxdxCxegedabagacad

PS1="\w \$ "
