# profile.d/dev-on-osx

LANG=en_US.UTF-8

export GIT_PAGER='less -FXR'

# (git) bash completion
if [ -f `brew --prefix`/etc/bash_completion ]; then
  . `brew --prefix`/etc/bash_completion
fi

# check the window size after each command and, if necessary,
# update the values of LINES and COLUMNS.
shopt -s checkwinsize

export CLICOLOR=1
export LSCOLORS=DxGxcxdxCxegedabagacad

PS1="\w \$ "

alias ma='git checkout master'
